# var/lib/dpkg/info/dropbear.list
/.
/usr
/usr/bin
/usr/bin/dropbearkey
/usr/bin/dbclient
/usr/share
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/dropbearkey.1.gz
/usr/share/man/man1/dbclient.1.gz
/usr/share/man/man1/dropbearconvert.1.gz
/usr/share/man/man8
/usr/share/man/man8/dropbear.8.gz
/usr/share/initramfs-tools
/usr/share/initramfs-tools/conf-hooks.d
/usr/share/initramfs-tools/conf-hooks.d/dropbear
/usr/share/initramfs-tools/scripts
/usr/share/initramfs-tools/scripts/init-bottom
/usr/share/initramfs-tools/scripts/init-bottom/dropbear
/usr/share/initramfs-tools/scripts/init-premount
/usr/share/initramfs-tools/scripts/init-premount/devpts
/usr/share/initramfs-tools/scripts/init-premount/dropbear
/usr/share/initramfs-tools/hooks
/usr/share/initramfs-tools/hooks/dropbear
/usr/share/doc
/usr/share/doc/dropbear
/usr/share/doc/dropbear/README
/usr/share/doc/dropbear/TODO
/usr/share/doc/dropbear/README.Debian
/usr/share/doc/dropbear/changelog.Debian.gz
/usr/share/doc/dropbear/changelog.gz
/usr/share/doc/dropbear/copyright
/usr/share/doc/dropbear/README.Debian.diet
/usr/share/doc/dropbear/README.runit
/usr/lib
/usr/lib/dropbear
/usr/lib/dropbear/dropbearconvert
/usr/sbin
/usr/sbin/dropbear
/etc
/etc/init.d
/etc/init.d/dropbear
/etc/dropbear
/etc/dropbear/run
/etc/dropbear/log
/etc/dropbear/log/run
/etc/dropbear/log/main
