# var/lib/dpkg/info/dropbear.conffiles
/etc/init.d/dropbear
/etc/dropbear/run
/etc/dropbear/log/run
/usr/share/initramfs-tools/conf-hooks.d/dropbear
