# var/lib/dpkg/info/mdadm.list
/.
/etc
/etc/logcheck
/etc/logcheck/violations.d
/etc/logcheck/violations.d/mdadm
/etc/logcheck/ignore.d.server
/etc/logcheck/ignore.d.server/mdadm
/etc/cron.d
/etc/cron.d/mdadm
/etc/cron.daily
/etc/cron.daily/mdadm
/etc/init.d
/etc/init.d/mdadm
/etc/init.d/mdadm-waitidle
/etc/init.d/mdadm-raid
/etc/modprobe.d
/etc/modprobe.d/mdadm.conf
/etc/mdadm
/sbin
/sbin/mdadm
/sbin/mdmon
/lib
/lib/systemd
/lib/systemd/system-shutdown
/lib/systemd/system-shutdown/mdadm.shutdown
/lib/systemd/system
/lib/systemd/system/mdadm-grow-continue@.service
/lib/systemd/system/mdadm-last-resort@.service
/lib/systemd/system/mdadm-last-resort@.timer
/lib/systemd/system/mdmonitor.service
/lib/systemd/system/mdmon@.service
/lib/udev
/lib/udev/rules.d
/lib/udev/rules.d/64-md-raid-assembly.rules
/lib/udev/rules.d/63-md-raid-arrays.rules
/usr
/usr/share
/usr/share/doc-base
/usr/share/doc-base/mdadm-readme-recipes
/usr/share/doc-base/mdadm-faq
/usr/share/doc
/usr/share/doc/mdadm
/usr/share/doc/mdadm/TODO.gz
/usr/share/doc/mdadm/README.recipes.gz
/usr/share/doc/mdadm/FAQ.gz
/usr/share/doc/mdadm/ANNOUNCE-3.2.4.gz
/usr/share/doc/mdadm/external-reshape-design.txt.gz
/usr/share/doc/mdadm/mdmon-design.txt.gz
/usr/share/doc/mdadm/README.Debian.gz
/usr/share/doc/mdadm/changelog.Debian.gz
/usr/share/doc/mdadm/NEWS.Debian.gz
/usr/share/doc/mdadm/changelog.gz
/usr/share/doc/mdadm/examples
/usr/share/doc/mdadm/examples/syslog-events
/usr/share/doc/mdadm/examples/mdadm.conf-example
/usr/share/doc/mdadm/copyright
/usr/share/doc/mdadm/TODO.Debian
/usr/share/doc/mdadm/ANNOUNCE-3.3.2
/usr/share/doc/mdadm/ANNOUNCE-3.3.1
/usr/share/doc/mdadm/ANNOUNCE-3.3
/usr/share/doc/mdadm/ANNOUNCE-3.2.6
/usr/share/doc/mdadm/ANNOUNCE-3.2.5
/usr/share/doc/mdadm/ANNOUNCE-3.2.3
/usr/share/doc/mdadm/ANNOUNCE-3.2.2
/usr/share/doc/mdadm/ANNOUNCE-3.2.1
/usr/share/doc/mdadm/ANNOUNCE-3.2
/usr/share/doc/mdadm/ANNOUNCE-3.1.5
/usr/share/doc/mdadm/ANNOUNCE-3.1.4
/usr/share/doc/mdadm/ANNOUNCE-3.1.3
/usr/share/doc/mdadm/ANNOUNCE-3.1.2
/usr/share/doc/mdadm/ANNOUNCE-3.1.1
/usr/share/doc/mdadm/ANNOUNCE-3.1
/usr/share/doc/mdadm/ANNOUNCE-3.0.3
/usr/share/doc/mdadm/ANNOUNCE-3.0.2
/usr/share/doc/mdadm/ANNOUNCE-3.0.1
/usr/share/doc/mdadm/ANNOUNCE-3.0
/usr/share/doc/mdadm/README.checkarray
/usr/share/bug
/usr/share/bug/mdadm
/usr/share/bug/mdadm/presubj
/usr/share/bug/mdadm/script
/usr/share/mdadm
/usr/share/mdadm/checkarray
/usr/share/mdadm/mkconf
/usr/share/initramfs-tools
/usr/share/initramfs-tools/scripts
/usr/share/initramfs-tools/scripts/local-top
/usr/share/initramfs-tools/scripts/local-top/mdadm
/usr/share/initramfs-tools/hooks
/usr/share/initramfs-tools/hooks/mdadm
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/mdadm.conf.5.gz
/usr/share/man/man4
/usr/share/man/man4/md.4.gz
/usr/share/man/man8
/usr/share/man/man8/mdmon.8.gz
/usr/share/man/man8/mdadm.8.gz
/lib/systemd/system/mdadm.service
/lib/systemd/system/mdadm-waitidle.service
